# ansible/roles/worker/templates/concourse_upstart.conf.j2
start on runlevel [2345]
stop on runlevel [06]

script
dtach -n {{ dtach_file }} {{ base_dir }}/bootstrap.sh
end script
